# supervisor.ini
[program:gvue-scaffold]
command=/yourpath/bin/server -f conf.prod.yml
directory=/yourpath
autostart=true
autorestart=true
startsecs=1
startretries=3
user=user
redirect_stderr=true
stdout_logfile=/var/logs/gvue-scaffold.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10